# recipes-bsp/tn-growpart-helper/tn-growpart-helper/tn-growpart-helper_xwindow.service
[Unit]
Description=Resize root filesystem to fit available disk space
After=systemd-remount-fs.service xserver-nodm.service
Wants=xserver-nodm.service
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
StandardOutput=journal+console
ExecStart=-/usr/sbin/tn-growpart-helper
ExecStartPost=/bin/systemctl disable tn-growpart-helper.service
ExecStop=/sbin/reboot

[Install]
WantedBy=multi-user.target
