# recipes-connectivity/bluez5/files/serial-qcabtfw@.timer
[Unit]
Description=Timer for QCA Bluetooth firmware service on %I

[Timer]
OnStartupSec=3

[Install]
WantedBy=timers.target
